# DevOps/homeworks/hw5/src/watcher.service
[Unit]
Description=Watch ~/watch directory and process new files
After=network.target

[Service]
Type=simple
User=slavko
WorkingDirectory=/home/slavko
ExecStart=/home/slavko/watcher.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
